# deploy/imcalc.service
[Unit]
Description=Import Calculator (Next.js)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/imcalc/app
Environment=NODE_ENV=production
EnvironmentFile=/var/www/imcalc/app/.env.local
ExecStart=/usr/bin/npm run start -- -p 3000 -H 127.0.0.1
Restart=on-failure
RestartSec=5

# Hardening (optional; ensure User can read .env.local and APP_DATA_DIR)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
